# deploy/systemd/backuper-server.service
[Unit]
Description=Backuper server (TLS dialer, sync initiator)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=backuper
Group=backuper
EnvironmentFile=-/etc/backuper/.env
WorkingDirectory=/var/lib/backuper
ExecStart=/usr/local/bin/backuper-server run
Restart=always
RestartSec=5
TimeoutStopSec=120
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
ReadWritePaths=/srv/backuper /var/lib/backuper /var/log/backuper

[Install]
WantedBy=multi-user.target
