# deploy/systemd/backuper-client.service
[Unit]
Description=Backuper client (TLS listener, backup source)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=backuper
Group=backuper
EnvironmentFile=-/etc/backuper/.env
WorkingDirectory=/var/log/backuper
ExecStart=/usr/local/bin/backuper-client run
Restart=always
RestartSec=5
TimeoutStopSec=120
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
# BACKUP_DIR должен быть в ReadWritePaths (для restore). По умолчанию /data.
ReadWritePaths=/data /var/log/backuper

[Install]
WantedBy=multi-user.target
